# dot_config/systemd/user/mailsync.timer
[Unit]
Description=mailsync

[Timer]
# Every 2 min
OnCalendar=*:0/2
Persistent=true
Unit=mailsync.service

[Install]
WantedBy=default.target
