# dot_config/systemd/user/mailsync.service.tmpl
[Unit]
Description=mailsync
Wants=mailsync.timer

[Service]
Type=oneshot
ExecStart=/home/{{ .chezmoi.username }}/bin/mailsync

[Install]
WantedBy=default.target
