# dot_config/systemd/user/sshfzf.service.tmpl
[Unit]
Description=Update sshfzf cache

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/home//{{ .chezmoi.username }}/bin/sshfzfupdatecache.sh

[Install]
WantedBy=default.target
