# dot_config/systemd/user/sshfzf.service
[Unit]
Description=Update sshfzf cache

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/home/kgn/bin/sshfzfupdatecache.sh

[Install]
WantedBy=default.target
