# dot_config/systemd/user/wallpaper.service.tmpl
[Unit]
Description=Wallpaper
Wants=wallpaper.timer

[Service]
Type=oneshot
ExecStart=/home/{{ .chezmoi.username }}/.config/hypr/scripts/changewallpaper.sh

[Install]
WantedBy=default.target
